# roles/autoup/templates/etc_systemd_system_autoup.timer.j2
[Unit]
Description=AutoUp Timer

[Timer]
OnBootSec=0min
OnUnitActiveSec=1min
Unit=autoup.service

[Install]
WantedBy=timers.target
